# systemd/otto-bgp-rpki-update.service
[Unit]
Description=Otto BGP RPKI Cache Update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=otto-bgp
Group=otto-bgp
WorkingDirectory=/var/lib/otto-bgp/rpki
ExecStart=/usr/bin/rpki-client -j -o /var/lib/otto-bgp/rpki/vrp_cache.json
StandardOutput=journal
StandardError=journal
NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
ReadWritePaths=/var/lib/otto-bgp/rpki

[Install]
WantedBy=multi-user.target
